# sparse/lib/systemd/system/sensorfwd.service.d/10-wait_for_sensors.conf
[Service]
ExecStartPre=/bin/sh -c "[[ $(/usr/bin/getprop init.svc.sensors) == "running" ]]"
RestartSec=2

